# deploy/supervisor-laravel-worker.ini
; Supervisor config for Laravel queue worker
;
; 1. Edit APP_PATH and user= below (use your real app path and Linux user).
; 2. Copy to Supervisor's config dir, e.g.:
;      sudo cp deploy/supervisor-laravel-worker.ini /etc/supervisord.d/laravel-worker.ini
; 3. sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start laravel-worker:*

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home2/royalce1/laravel-app/school-management-system/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=royalce1
numprocs=1
redirect_stderr=true
stdout_logfile=/home2/royalce1/laravel-app/school-management-system/storage/logs/worker.log
stopwaitsecs=60
