# config/supervisor-erp-worker.conf
; Supervisor config for Laravel queue worker (EC2/server)
; Place in /etc/supervisor/conf.d/erp-worker.conf on the server
; Then: sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start erp-worker:*

[program:erp-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/erp/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/erp/storage/logs/worker.log
stopwaitsecs=3600
